# root/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv sh

UMASK_SET=${UMASK_SET:-022}
WEBUI_PORT=${WEBUI_PORT:-9000}

umask "$UMASK_SET"

exec \
	s6-setuidgid user /usr/local/bin/qbittorrent-nox --webui-port="${WEBUI_PORT}"